# my-personal-PKGBUILDs/checkupdates-userservice/checkupdates.install
post_install() {
	echo "****************************************************************************************"
	echo "Run 'systemctl --user enable --now checkupdates.timer' to autostart this service on boot"
	echo "****************************************************************************************"
}
