# my-personal-PKGBUILDs/sshagent-userservice/sshagent.install
post_install() {
	echo "**************************************************************************************"
	echo "Run 'systemctl --user enable --now sshagent.service' to autostart this service on boot"
	echo "**************************************************************************************"
}
